# rpm_buildtree/nDeploy-cluster-slave-pkg-centos7/usr/lib/systemd/system/ndeploy_hhvm@.service
[Service]
ExecStart=/usr/local/bin/hhvm --config /opt/nDeploy/hhvm.d/%i.ini --user %i --mode daemon

[Install]
WantedBy=multi-user.target
