# rpm_buildtree/nDeploy-pkg/etc/rc.d/init.d/ndeploy_backends
#!/bin/bash
#
# chkconfig: 35 90 12
# description: nDeploy Backends init and stop script
#
# Get function from functions library
. /etc/init.d/functions
# Start
start() {
	/opt/nDeploy/scripts/init_backends.py start
}
# stop
stop() {
	/opt/nDeploy/scripts/init_backends.py stop
}
#Conf reload
reload() {
	/opt/nDeploy/scripts/init_backends.py reload
}
### main logic ###
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  reload)
  		reload
		;;
  restart)
        stop
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|reload|restart}"
        exit 1
esac
exit 0
