# rpm_buildtree/nDeploy-pkg-centos7/etc/cron.d/ndeploy
* * * * * root /opt/nDeploy/scripts/config_generation_queue_runner.sh
